# worker/crontab
0 * * * * node cron-updateDBApod.js
0 6 * * * node cron-sendMediaSubscribers.js
